# bin/rheaStartMITM.service
[Unit]
Description=Fa partire il MITM (Man In The Middle)
After=multi-user.target

[Service]
WorkingDirectory=/home/pi/rhea
User=pi
ExecStart=/home/pi/rhea/RASPI_RELEASE_rasPI
#ExecStart=/home/pi/rhea/RASPI_DEBUG_rasPI

[Install]
WantedBy=multi-user.target

