# bin/rheaStartRasPIESAPI.service
[Unit]
Description=Start rhea bootloader for rasPI prg
After=network.target

[Service]
WorkingDirectory=/home/pi/rhea/gpu-fts-nestle-2019/bin
User=pi
ExecStart=/home/pi/rhea/gpu-fts-nestle-2019/bin/RASPI_RELEASE_rasPIBootLoader

[Install]
WantedBy=network.target

